# deploy/voice-assistant-portal.service.template
[Unit]
Description=Jetson Voice Assistant Admin Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__JETSON_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=__APP_DIR__/venv/bin/python __APP_DIR__/src/admin_portal.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
